# nad-rs232-rest.service
[Unit]
Description=NAD RS232 REST-API Service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python /usr/local/bin/nad-rs232-rest.py
Restart=always

[Install]
WantedBy=multi-user.target
